# .config/systemd/toggle-leds.service
[Unit]
Description=toggles off annoying thinkpad leds

[Service]
ExecStart=/home/leah/.local/bin/thinkpad-leds

[Install]
WantedBy=multi-user.target